# scripts/user-backup.timer
[Unit]
Description=Run AI Orchestrator backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
